# logrotate/config/crontab
15 14 * * *     /usr/sbin/logrotate -v /tokaido/config/logrotate/logrotate.conf