# logrotate/config/crontab
0 0 * * *     /usr/sbin/logrotate -v /app/config/logrotate/logrotate.conf